# snapshots/snapshots.cron
SHELL=/bin/bash
PATH=/sbin:/bin:/usr/sbin:/usr/bin
MAILTO=root
30 0,3,6,9,12,15,18,21 * * * root snapshots.sh -ns
